# adblock_sync_nftset
#!/bin/sh /etc/rc.common
START=99
STOP=15

# --- 配置区 ---
URL="https://anti-ad.net/domains.txt"
MEM_DIR="/tmp/dnsmasq_ads"
MEM_CONF="$MEM_DIR/ad_processed.conf"
CUSTOM_CONF_DIR="/root/.dnsmasq_custom"
BACKUP_DIR="$MEM_DIR/backups"
MAX_BACKUPS=1
LOCKFILE="/var/run/adblock_sync.lock"

log() {
    logger -t AdBlock "$1"
}

fail() {
    log "错误：$1"
    exit 1
}

network_ready_or_wait() {
    local count=0
    local targets="223.5.5.5 www.baidu.com"
    log "正在检测网络连接..."
    while true; do
        for target in $targets; do
            if ping -c 1 -W 2 $target >/dev/null 2>&1; then
                log "网络已就绪。"
                return 0
            fi
        done
        count=$((count + 1))
        if [ $count -gt 30 ]; then
            log "网络连接失败，停止本次更新。"
            exit 0
        fi
        sleep 2
    done
}

sync_time() {
    if command -v ntpd >/dev/null; then
        # 从 UCI 获取 NTP 服务器列表
        local ntp_servers=$(uci get system.ntp.server 2>/dev/null)
        if [ -z "$ntp_servers" ]; then
            log "未从 UCI 获取到 NTP 服务器，使用默认 ntp.aliyun.com"
            ntp_servers="ntp.aliyun.com"
        fi
        # 构造参数
        local ntp_args=""
        for s in $ntp_servers; do
            ntp_args="$ntp_args -p $s"
        done
        log "正在同步时间 (ntpd -n -q$ntp_args)..."
        ntpd -n -q $ntp_args 2>&1 | logger -t AdBlock
    else
        log "警告：未找到时间同步工具，跳过时间同步"
        return 1
    fi
}

check_dnsmasq_nftset_support() {
    log "检查 dnsmasq 是否支持 nftset..."
    echo 'nftset=/example.com/4#inet#fw4#test' | dnsmasq --test --conf-file=/dev/stdin 2>/dev/null
    if [ $? -ne 0 ]; then
        fail "当前 dnsmasq 不支持 nftset，请升级为 dnsmasq-full"
    fi
    log "dnsmasq 支持 nftset，继续..."
}

setup_firewall() {
    local changed=0
    local items="
ipset|ad_block_v4|ad_block_v4|ipv4
ipset|ad_block_v6|ad_block_v6|ipv6
rule|ad_block_v4_rule|ad_block_v4|ipv4
rule|ad_block_v6_rule|ad_block_v6|ipv6
"
    local IFS='
'
    for item in $items; do
        [ -z "$item" ] && continue
        local type=$(echo "$item" | cut -d'|' -f1)
        local conf_name=$(echo "$item" | cut -d'|' -f2)
        local ipset_name=$(echo "$item" | cut -d'|' -f3)
        local family=$(echo "$item" | cut -d'|' -f4)

        local exists=""
        local sections=$(uci show firewall | grep -o "firewall\.@${type}\[[0-9]\+\]" | sort -u)
        for sec in $sections; do
            local name=$(uci -q get "${sec}.name")
            if [ "$name" = "$conf_name" ]; then
                exists="$sec"
                break
            fi
        done

        if [ -z "$exists" ]; then
            log "添加 $type: $conf_name (ipset=$ipset_name, family=$family)"
            local section=$(uci add firewall "$type")
            if [ "$type" = "ipset" ]; then
                uci set firewall."$section".name="$ipset_name"
                uci set firewall."$section".family="$family"
                uci add_list firewall."$section".match="dest_net"
            else
                uci set firewall."$section".name="$conf_name"
                uci set firewall."$section".src="lan"
                uci set firewall."$section".dest="wan"
                uci set firewall."$section".family="$family"
                uci set firewall."$section".ipset="$ipset_name"
                uci set firewall."$section".target="REJECT"
            fi
            changed=1
        else
            log "跳过 $type: $conf_name (已存在)"
        fi
    done

    if [ $changed -eq 1 ]; then
        uci commit firewall
        if command -v fw4 >/dev/null; then
            fw4 reload >/dev/null 2>&1 && log "防火墙已重载" || log "警告：防火墙重载失败"
        else
            log "警告：fw4 命令不存在，请手动重载防火墙"
        fi
    else
        log "防火墙配置已存在，无需修改"
    fi
}

acquire_lock() {
    if ! mkdir "$LOCKFILE" 2>/dev/null; then
        log "另一个更新实例正在运行，退出。"
        exit 1
    fi
}

release_lock() {
    rmdir "$LOCKFILE" 2>/dev/null || true
}

prepare_environment() {
    mkdir -p "$MEM_DIR"
    if [ -d "$CUSTOM_CONF_DIR" ]; then
        cp "$CUSTOM_CONF_DIR"/*.conf "$MEM_DIR/" 2>/dev/null
        log "已同步本地自定义配置至内存。"
    fi

    if ! uci -q get dhcp.@dnsmasq[0].confdir | grep -Fx "$MEM_DIR" >/dev/null; then
        uci add_list dhcp.@dnsmasq[0].confdir="$MEM_DIR"
        uci commit dhcp
        log "已成功将 $MEM_DIR 写入 dnsmasq 配置。"
    fi
}

download_and_process() {
    local download_cmd
    if command -v curl >/dev/null 2>&1; then
        download_cmd="curl -4 --connect-timeout 10 --max-time 120 -sL"
    elif command -v wget >/dev/null 2>&1; then
        download_cmd="wget --no-check-certificate -4 -q --timeout=120 -O -"
    else
        fail "未找到 curl 或 wget，无法下载"
    fi

    log "开始下载广告域名列表..."
    if ! $download_cmd "$URL" | awk '{
        gsub(/\r/,"");
        if ($1 !~ /^#/ && $1 != "") {
            printf "nftset=/%s/4#inet#fw4#ad_block_v4,6#inet#fw4#ad_block_v6\n", $1
        }
    }' | sort -u > "$TMP_DIR/ad_processed.conf.tmp"; then
        fail "下载失败，请检查网络或URL"
    fi

    if [ ! -s "$TMP_DIR/ad_processed.conf.tmp" ]; then
        fail "下载内容为空"
    fi

    # 语法校验
    dnsmasq --test --conf-file="$TMP_DIR/ad_processed.conf.tmp" || fail "语法校验失败"
}

backup_current_conf() {
    mkdir -p "$BACKUP_DIR"
    if [ -f "$MEM_CONF" ]; then
        BACKUP_NAME="$BACKUP_DIR/ad_$(date +%Y%m%d-%H%M%S).conf"
        cp -p "$MEM_CONF" "$BACKUP_NAME" || fail "备份失败"
        ls -t "$BACKUP_DIR"/ad_*.conf 2>/dev/null | tail -n +$((MAX_BACKUPS+1)) | xargs rm -f
    else
        BACKUP_NAME=""
    fi
}

flush_nftsets() {
    if command -v nft >/dev/null; then
        log "正在清空 nftset 集合..."
        if nft list set inet fw4 ad_block_v4 >/dev/null 2>&1; then
            nft flush set inet fw4 ad_block_v4 && log "已清空集合 inet fw4 ad_block_v4" || log "警告：清空集合 inet fw4 ad_block_v4 失败"
        else
            log "集合 inet fw4 ad_block_v4 不存在，跳过"
        fi
        if nft list set inet fw4 ad_block_v6 >/dev/null 2>&1; then
            nft flush set inet fw4 ad_block_v6 && log "已清空集合 inet fw4 ad_block_v6" || log "警告：清空集合 inet fw4 ad_block_v6 失败"
        else
            log "集合 inet fw4 ad_block_v6 不存在，跳过"
        fi
    else
        log "警告：nft 命令未找到，跳过清空集合操作"
    fi
}

apply_new_config() {
    mv -f "$TMP_DIR/ad_processed.conf.tmp" "$MEM_CONF" || fail "替换文件失败"
    lines=$(wc -l < "$MEM_CONF")
    log "生成 $lines 条规则，正在重启 dnsmasq..."

    if ! /etc/init.d/dnsmasq restart >/dev/null 2>&1; then
        # 重启失败，尝试回滚
        log "dnsmasq 重启失败，尝试回滚"
        if [ -n "$BACKUP_NAME" ] && [ -f "$BACKUP_NAME" ]; then
            cp -f "$BACKUP_NAME" "$MEM_CONF"
            /etc/init.d/dnsmasq restart >/dev/null 2>&1
        fi
        fail "dnsmasq 重启失败，已回滚"
    fi

    sleep 4
    local try=0
    local max_tries=6
    local success=0
    while [ $try -lt $max_tries ]; do
        if nslookup www.baidu.com 127.0.0.1 >/dev/null 2>&1; then
            success=1
            break
        fi
        try=$((try+1))
        [ $try -lt $max_tries ] && sleep 2
    done

    if [ $success -eq 1 ]; then
        log "更新完成，拦截已生效"
    else
        log "警告：dnsmasq 启动后无法正常解析域名，执行回滚"
        if [ -n "$BACKUP_NAME" ] && [ -f "$BACKUP_NAME" ]; then
            cp -f "$BACKUP_NAME" "$MEM_CONF"
            /etc/init.d/dnsmasq restart >/dev/null 2>&1
        fi
        fail "回滚完成，请检查原因"
    fi
}

update_list() {
    network_ready_or_wait
    sync_time

    acquire_lock
    TMP_DIR=$(mktemp -d) || fail "无法创建临时目录"
    trap "rm -rf '$TMP_DIR' 2>/dev/null; rmdir '$LOCKFILE' 2>/dev/null" EXIT INT TERM

    check_dnsmasq_nftset_support
    setup_firewall
    prepare_environment
    download_and_process
    backup_current_conf
    flush_nftsets
    apply_new_config

    release_lock
    rm -rf "$TMP_DIR"

    log "广告列表更新进程完成"
}

start() {
    (sleep 60; update_list) &
}

reload() {
    update_list
}

stop() {
    :
}

restart() {
    reload
}
