# adblock_sync_nftset
#!/bin/sh /etc/rc.common

START=99
STOP=15

# --- 配置区 ---
URL="https://anti-ad.net/domains.txt"
MEM_DIR="/tmp/dnsmasq_ads"
MEM_CONF="$MEM_DIR/ad_processed.conf"
CUSTOM_CONF_DIR="/root/.dnsmasq_custom"
BACKUP_DIR="$MEM_DIR/backups"
MAX_BACKUPS=1
LOCKFILE="/var/run/adblock_sync.lock"

log() {
    logger -t AdBlock "$1"
}

fail() {
    log "错误：$1"
    exit 1
}

sync_time() {
    local ntp_server="ntp.aliyun.com"
    if command -v ntpd >/dev/null; then
        log "正在同步时间 (ntpd -q)..."
        ntpd -q -p "$ntp_server" 2>&1 | logger -t AdBlock
    elif command -v ntpdate >/dev/null; then
        log "正在同步时间 (ntpdate)..."
        ntpdate "$ntp_server" 2>&1 | logger -t AdBlock
    elif command -v rdate >/dev/null; then
        log "正在同步时间 (rdate)..."
        rdate -s "$ntp_server" 2>&1 | logger -t AdBlock
    else
        log "警告：未找到时间同步工具，跳过时间同步"
        return 1
    fi
}

setup_firewall() {
    local changed=0
    local items="
ipset|ad_block_v4|ad_block_v4|ipv4
ipset|ad_block_v6|ad_block_v6|ipv6
rule|ad_block_v4_rule|ad_block_v4|ipv4
rule|ad_block_v6_rule|ad_block_v6|ipv6
"
    local IFS='
'
    for item in $items; do
        [ -z "$item" ] && continue
        local type=$(echo "$item" | cut -d'|' -f1)
        local conf_name=$(echo "$item" | cut -d'|' -f2)
        local ipset_name=$(echo "$item" | cut -d'|' -f3)
        local family=$(echo "$item" | cut -d'|' -f4)

        local exists=""
        local sections=$(uci show firewall | grep -o "firewall\.@${type}\[[0-9]\+\]" | sort -u)
        for sec in $sections; do
            local name=$(uci -q get "${sec}.name")
            if [ "$name" = "$conf_name" ]; then
                exists="$sec"
                break
            fi
        done

        if [ -z "$exists" ]; then
            log "添加 $type: $conf_name (ipset=$ipset_name, family=$family)"
            local section=$(uci add firewall "$type")
            if [ "$type" = "ipset" ]; then
                uci set firewall."$section".name="$ipset_name"
                uci set firewall."$section".family="$family"
                uci add_list firewall."$section".match="dest_net"
            else
                uci set firewall."$section".name="$conf_name"
                uci set firewall."$section".src="lan"
                uci set firewall."$section".dest="wan"
                uci set firewall."$section".family="$family"
                uci set firewall."$section".ipset="$ipset_name"
                uci set firewall."$section".target="REJECT"
            fi
            changed=1
        else
            log "跳过 $type: $conf_name (已存在)"
        fi
    done

    if [ $changed -eq 1 ]; then
        uci commit firewall
        if command -v fw4 >/dev/null; then
            fw4 reload >/dev/null 2>&1 && log "防火墙已重载" || log "警告：防火墙重载失败"
        else
            log "警告：fw4 命令不存在，请手动重载防火墙"
        fi
    else
        log "防火墙配置已存在，无需修改"
    fi
}

update_list() {
    local count=0
    local targets="223.5.5.5 www.baidu.com"
    log "正在检测网络连接..."
    while true; do
        for target in $targets; do
            if ping -c 1 -W 2 $target >/dev/null 2>&1; then
                log "网络已就绪。"
                break 2
            fi
        done
        count=$((count + 1))
        if [ $count -gt 30 ]; then
            log "网络连接失败，停止本次更新。"
            exit 0 
        fi
        sleep 2
    done

    sync_time

    if ! mkdir "$LOCKFILE" 2>/dev/null; then
        log "另一个更新实例正在运行，退出。"
        exit 1
    fi
    trap "rmdir '$LOCKFILE' 2>/dev/null" EXIT INT TERM

    TMP_DIR=$(mktemp -d) || fail "无法创建临时目录"
    trap "rm -rf '$TMP_DIR' 2>/dev/null; rmdir '$LOCKFILE' 2>/dev/null" EXIT INT TERM

    log "检查 dnsmasq 是否支持 nftset..."
    echo 'nftset=/example.com/4#inet#fw4#test' | dnsmasq --test --conf-file=/dev/stdin 2>/dev/null
    if [ $? -ne 0 ]; then
        fail "当前 dnsmasq 不支持 nftset，请升级为 dnsmasq-full"
    fi
    log "dnsmasq 支持 nftset，继续..."

    setup_firewall

    mkdir -p "$MEM_DIR"
    if [ -d "$CUSTOM_CONF_DIR" ]; then
        cp "$CUSTOM_CONF_DIR"/*.conf "$MEM_DIR/" 2>/dev/null
        log "已同步本地自定义配置至内存。"
    fi

    if ! uci -q get dhcp.@dnsmasq[0].confdir | grep -Fx "$MEM_DIR" >/dev/null; then
        uci add_list dhcp.@dnsmasq[0].confdir="$MEM_DIR"
        uci commit dhcp
        log "已成功将 $MEM_DIR 写入 dnsmasq 配置。"
    fi

    log "开始下载广告域名列表..."
    local download_cmd download_arg download_prog
    if command -v curl >/dev/null; then
        download_cmd="curl -4 --connect-timeout 10 --max-time 120 -sL -o"
        download_arg="$TMP_DIR/raw_domains.txt"
        download_prog="curl"
    elif command -v wget >/dev/null; then
        download_cmd="wget --no-check-certificate -4 --timeout=120 --tries=1 -q -O"
        download_arg="$TMP_DIR/raw_domains.txt"
        download_prog="wget"
    else
        fail "未找到 curl 或 wget，无法下载"
    fi

    log "使用 $download_prog 下载列表..."
    if ! $download_cmd "$download_arg" "$URL"; then
        fail "下载失败（$download_prog 返回非零状态）"
    fi
    if [ ! -s "$TMP_DIR/raw_domains.txt" ]; then
        fail "下载内容为空"
    fi

    # 处理下载的原始文件，生成 dnsmasq 规则
    awk '{
        gsub(/\r/,""); 
        if ($1 !~ /^#/ && $1 != "") {
            printf "nftset=/%s/4#inet#fw4#ad_block_v4,6#inet#fw4#ad_block_v6\n", $1
        }
    }' "$TMP_DIR/raw_domains.txt" | sort -u > "$TMP_DIR/ad_processed.conf.tmp" || fail "处理下载文件失败"

    rm -f "$TMP_DIR/raw_domains.txt"

    dnsmasq --test --conf-file="$TMP_DIR/ad_processed.conf.tmp" || fail "语法校验失败"

    mkdir -p "$BACKUP_DIR"
    if [ -f "$MEM_CONF" ]; then
        BACKUP_NAME="$BACKUP_DIR/ad_$(date +%Y%m%d-%H%M%S).conf"
        cp -p "$MEM_CONF" "$BACKUP_NAME" || fail "备份失败"
        ls -t "$BACKUP_DIR"/ad_*.conf 2>/dev/null | tail -n +$((MAX_BACKUPS+1)) | xargs rm -f
    fi

    if command -v nft >/dev/null; then
        log "正在清空 nftset 集合..."
        if nft list set inet fw4 ad_block_v4 >/dev/null 2>&1; then
            nft flush set inet fw4 ad_block_v4 && log "已清空集合 inet fw4 ad_block_v4" || log "警告：清空集合 inet fw4 ad_block_v4 失败"
        else
            log "集合 inet fw4 ad_block_v4 不存在，跳过"
        fi
        if nft list set inet fw4 ad_block_v6 >/dev/null 2>&1; then
            nft flush set inet fw4 ad_block_v6 && log "已清空集合 inet fw4 ad_block_v6" || log "警告：清空集合 inet fw4 ad_block_v6 失败"
        else
            log "集合 inet fw4 ad_block_v6 不存在，跳过"
        fi
    else
        log "警告：nft 命令未找到，跳过清空集合操作"
    fi

    mv -f "$TMP_DIR/ad_processed.conf.tmp" "$MEM_CONF" || fail "替换文件失败"
    lines=$(wc -l < "$MEM_CONF")
    log "生成 $lines 条规则，正在重启 dnsmasq..."

    if /etc/init.d/dnsmasq restart >/dev/null 2>&1; then
        sleep 4
        local try=0
        local max_tries=6
        local success=0
        while [ $try -lt $max_tries ]; do
            if nslookup www.baidu.com 127.0.0.1 >/dev/null 2>&1; then
                success=1
                break
            fi
            try=$((try+1))
            [ $try -lt $max_tries ] && sleep 2
        done

        if [ $success -eq 1 ]; then
            log "更新完成，拦截已生效"
        else
            log "警告：dnsmasq 启动后无法正常解析域名，执行回滚"
            [ -f "$BACKUP_NAME" ] && cp -f "$BACKUP_NAME" "$MEM_CONF"
            /etc/init.d/dnsmasq restart >/dev/null 2>&1
            fail "回滚完成，请检查原因"
        fi
    else
        fail "dnsmasq 重启失败"
    fi
}

start() {
    update_list
}

reload() {
    update_list
}

stop() {
    :
}

restart() {
    reload
}
